# systemd/internet-speed.service
[Unit]
Description=Internet Speed Monitoring Service

[Service]
Type=exec
ExecStart=python3 

[Install]
WantedBy=multi-user.target